# scripts/launchagents/com.luke.vos.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luke.vos.morning-brief</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/python@3.13/bin/python3</string>
        <string>/Volumes/MontereyT7/venture-os/components/morning-briefer/briefer.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
