# scripts/launchagents/com.luke.vos.routing-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luke.vos.routing-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/python@3.13/bin/python3</string>
        <string>/Volumes/MontereyT7/venture-os/components/routing-refresh/refresher.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/macbook/.claude/logs/routing-refresh.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbook/.claude/logs/routing-refresh.err</string>
</dict>
</plist>
